# deploy/claude-skills-supervisor.conf
# API 服务
[program:claude-skills-api]
command=/var/www/claude-skills-miniapp/backend/venv/bin/python -m uvicorn src.api.main:app --host 0.0.0.0 --port 8000
directory=/var/www/claude-skills-miniapp/backend
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/claude-skills-api.log
environment=PATH="/var/www/claude-skills-miniapp/backend/venv/bin"

# 内容更新任务 - 持续运行，自动断点续传
[program:claude-skills-content-update]
command=/var/www/claude-skills-miniapp/backend/venv/bin/python scripts/continuous_update_content.py
directory=/var/www/claude-skills-miniapp/backend
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/claude-skills-content-update.log
stderr_logfile=/var/log/claude-skills-content-update-error.log
environment=PATH="/var/www/claude-skills-miniapp/backend/venv/bin"
# 如果程序正常退出（所有内容更新完成），不要重启
autorestart=unexpected
exitcodes=0
# 启动延迟，避免和 API 服务抢资源
startsecs=10
# 停止等待时间
stopwaitsecs=30

# 翻译任务 - 定时运行（每天凌晨 2 点）
# 注意：这个不使用 supervisor，而是使用 crontab
# 见下方 crontab 配置说明
