# deploy/claude-skills-api.conf
[program:claude-skills-api]
command=/var/www/claude-skills-miniapp/backend/venv/bin/python -m uvicorn src.api.main:app --host 0.0.0.0 --port 8000
directory=/var/www/claude-skills-miniapp/backend
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/claude-skills-api.log
environment=PATH="/var/www/claude-skills-miniapp/backend/venv/bin"
